# egpusetup
#!/sbin/runscript
# Facts about this script:
# This was written over the course of 10 minutes (in fact, I spent more time writing these few lines)
# It requires you to blacklist all gpu drivers in the modprobe config to prevent autoload
# It will probably not work for you without modification (If you have a ThinkPad X220 it just might)
# I can and will take NO responsibility for any damage this may cause
# I will take full responsibility for anything awesome this may cause ;)

depend() {
	before xdm
}

start() {
	/usr/sbin/lspci | grep -i nvidia &> /dev/null
	if [ $? -eq 0 ]; then
		echo 1 > /sys/bus/pci/devices/0000:00:02.0/remove
		/sbin/modprobe nvidia
		/usr/bin/eselect opengl set nvidia
		/bin/cp /root/00-nvidia.conf /usr/share/X11/xorg.conf.d/
	else
		if [ -e /usr/share/X11/xorg.conf.d/00-nvidia.conf ]; then
			/bin/rm /usr/share/X11/xorg.conf.d/00-nvidia.conf
		fi
		/sbin/modprobe i915
		/usr/bin/eselect opengl set xorg-x11
	fi
}

stop() {
	if [ -e /usr/share/X11/xorg.conf.d/00-nvidia.conf ]; then
		/bin/rm /usr/share/X11/xorg.conf.d/00-nvidia.conf
	fi
}
